# configs/overlay_rootfs/var/spool/cron/crontabs/root
*/15 * * * * /usr/sbin/logrotate /etc/logrotate.conf
0 * * * * /scripts/remove_old.sh
0 2 * * 7 /scripts/reboot.sh
